# com.eduardocortez.systemoptimizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eduardocortez.systemoptimizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/eduardocortez/genius-code-1.0.1/system-optimizer/target/release/apollo-optimizer</string>
        <string>daemon</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/eduardocortez/Library/Logs/system-optimizer.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/eduardocortez/Library/Logs/system-optimizer.err.log</string>
</dict>
</plist>
