# scripts/com.eduardocortez.apollo-regression-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eduardocortez.apollo-regression-probe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/libexec/apollo-regression-probe.py</string>
        <string>--quiet</string>
    </array>
    <!-- 30 min: catches runtime regressions (thrashing/refault/boost-loop/
         sysctl-tamper) fast and builds a useful cycle-over-cycle trend in
         /var/lib/apollo/regression-metrics.jsonl without being noisy. -->
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/lib/apollo/regression-cron.log</string>
    <key>StandardErrorPath</key>
    <string>/var/lib/apollo/regression-cron.err</string>
</dict>
</plist>
