# scripts/com.eduardocortez.apollo-retrain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eduardocortez.apollo-retrain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/apollo-retrain.sh</string>
    </array>

    <!-- Run nightly at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Low priority: nice 15, throttle I/O -->
    <key>Nice</key>
    <integer>15</integer>
    <key>LowPriorityIO</key>
    <true/>
    <key>LowPriorityBackgroundIO</key>
    <true/>

    <!-- Resource limits: cap memory at 3GB to leave room for Apollo + OS -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>MemoryLock</key>
        <integer>3221225472</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>MemoryLock</key>
        <integer>3221225472</integer>
    </dict>

    <!-- Timeout: kill if running more than 30 minutes -->
    <key>TimeOut</key>
    <integer>1800</integer>

    <!-- Don't restart on failure -->
    <key>KeepAlive</key>
    <false/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/var/lib/apollo/retrain-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/var/lib/apollo/retrain-stderr.log</string>
</dict>
</plist>
